# infra/fairtrail-scrape.service
[Unit]
Description=Fairtrail flight price scrape
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
User=sotto
EnvironmentFile=/home/sotto/.config/fairtrail/env
ExecStart=/bin/bash -c 'curl -sf --max-time 600 -H "Authorization: Bearer ${CRON_SECRET}" http://localhost:3002/api/cron/scrape'
TimeoutStartSec=900
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
